# web-agent/armorcode-agent.service
[Unit]
Description=ArmorCode Web Agent
Documentation=https://github.com/armor-code/agent
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=500
StartLimitBurst=5

[Service]
Type=simple
ExecStart=/usr/bin/python3 worker.py --serverUrl=https://web-agent.armorcode.com --apiKey=<API_KEY>
Restart=on-failure
RestartSec=5s
User=root
WorkingDirectory=/opt/armorcode

# Environment variables
Environment=PYTHONPATH=/opt/armorcode
Environment=PYTHONUNBUFFERED=1

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=armorcode-agent

[Install]
WantedBy=multi-user.target